# deploy/systemd/btc-header-sync.timer
[Unit]
Description=Run BTC header relay sync every 10 minutes
Documentation=https://github.com/iriumlabs/irium/tree/main/deploy/systemd

[Timer]
# Wait until network is up and iriumd has had ~2 min to settle on boot
# before the first cycle, then fire 10 minutes after each completion.
OnBootSec=2min
OnUnitInactiveSec=10min
Unit=btc-header-sync.service
# Spread accidental synchronised wakeups across deployments by up to 30s.
RandomizedDelaySec=30s
# Run a missed tick on boot if the system was off when the timer should
# have fired.
Persistent=true

[Install]
WantedBy=timers.target
