# deploy/systemd/btc-header-sync.service
[Unit]
Description=Iriumd BTC header relay sync (one-shot, timer-driven)
Documentation=https://github.com/iriumlabs/irium/tree/main/deploy/systemd
After=network-online.target iriumd.service
Wants=network-online.target

[Service]
Type=oneshot
User=irium
Group=irium
EnvironmentFile=/etc/irium/btc-header-sync.env
ExecStart=/usr/local/bin/btc-header-sync
# stdout/stderr are mirrored to the log file by the binary itself; this
# extra hookup makes journalctl -u btc-header-sync.service also useful.
StandardOutput=journal
StandardError=journal
# Tighten capabilities — this binary only opens outbound HTTPS and
# appends to a single log file owned by the irium user.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/irium
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
