# systemd/iriumd.service
[Unit]
Description=Irium Rust Node
After=network-online.target
Wants=network-online.target

[Service]
User=@IRIUM_USER@
EnvironmentFile=-/etc/irium/iriumd.env
WorkingDirectory=@IRIUM_HOME@
ExecStart=@IRIUM_HOME@/target/release/iriumd
Restart=on-failure
RestartSec=5
TimeoutStopSec=90
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
