# systemd/iriumd.service.d/40-wants-stratums.conf
[Unit]
# When iriumd starts (cold start, post-stop, post-crash), pull in
# the four pool stratum services. PartOf in each stratum drop-in
# handles stop/restart cascade in the other direction.
# Fixes the post-stop+start scenario where PartOf alone leaves
# stratums inactive after iriumd is started from stopped state.
# (Issue: systemctl stop iriumd cascades to stop stratums via PartOf,
# but systemctl start iriumd does not cascade to start inactive stratums.)
Wants=irium-stratum.service irium-stratum-443.service irium-stratum-legacy.service irium-stratum-solo.service
